# files/etc/crontabs/root
*/5 * * * * /usr/local/bin/fccid
*/1 * * * * /usr/local/bin/rssi_monitor
