# files/etc/crontabs/root
0 */12 * * * /usr/local/bin/aredn_message.sh

