# files/etc/crontabs/root
*/5 * * * * /usr/local/bin/fccid
* * * * * /usr/local/bin/rssi_monitor
* * * * * /usr/local/bin/snrlog
* * * * * /usr/local/bin/clean_zombie.sh

