# files/etc/crontabs/root
*/5 * * * * /usr/local/bin/fccid
* * * * * /usr/local/bin/snrlog
