# files/etc/uci-defaults/10_dmz_mode_migrate
#! /bin/sh
if [ -f /etc/config/dmz-mode ] ; then
    mv /etc/config/dmz-mode /etc/dmz-mode
fi
