# files/etc/uci-defaults/13_opkg_useragent
#! /bin/sh
VERSION="$(cat /etc/mesh-release)"
echo "option user_agent opkg-${VERSION}" >> /etc/opkg.conf
