# files/etc/uci-defaults/93_lqm
#!/bin/sh

/sbin/uci -q set aredn.@lqm[0].enable=0
/sbin/uci -q commit aredn
/sbin/uci -c /etc/config.mesh -q set aredn.@lqm[0].enable=0
/sbin/uci -c /etc/config.mesh -q commit aredn
