# files/etc/uci-defaults/94_fix_wpad
#! /bin/sh
# Disable hostapd (etc) on devices without wifi

if [ ! -d /sys/class/ieee80211 ]; then
    /etc/init.d/wpad disable
fi
