# files/etc/uci-defaults/94_fix_wpad
#! /bin/sh
# Enable hostapd (etc) on devices with wifi in case it was disabled

if [ -d /sys/class/ieee80211 ]; then
    /etc/init.d/wpad enable
fi
