# files/etc/uci-defaults/94_update_permpkg
#! /bin/sh
# Update the permanent package list so we can't uninstall any of the standard system

opkg list-installed | sed s/\ .*// > /etc/permpkg
