# files/etc/uci-defaults/97_aredn_setup_meshstatus
#! /bin/sh
if [ "$(/sbin/uci -q get aredn.\@meshstatus[0])" = "" ]; then
    /sbin/uci -q add aredn meshstatus
    /sbin/uci -q commit aredn
fi
