# files/etc/init.d/local
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=99
boot() {
  
  # extract kamikaze generated vlan config
  #CMLARA Not needed on BULLETM's  maybe needed on other devices
  #[ -f /etc/vlan ] || head -n5 /etc/config/network > /etc/vlan
  
  # setup nvram variables
  [ -x /usr/local/bin/nvram-setup ] &&  /usr/local/bin/nvram-setup

  # run mode specific setup
  [ -x /etc/config/local ] && /etc/config/local
  [ -x /etc/local/services ] && /etc/local/services
}
