# files/etc/init.d/local-early
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=10
boot() {
    mkdir -p /tmp/dnsmasq.d
    mkdir -p /tmp/dnshosts.d
}
