# files/etc/hotplug.d/ntp/10-ntp-stratum
STATEFILE="/var/state/ntp-stratum"

if [[ "${stratum}" != "" ]]; then
	echo "${stratum}" > "${STATEFILE}"
fi
