# files/etc/xinetd.d/dotdraw
# description: a workaround for the overly restrictive olsr dot_draw plugin
# - open up the port to anyone instead of just one specific ip address

service dotdraw
{
        disable         = no
        socket_type     = stream
        protocol        = tcp
        wait            = no
        user            = root
        group           = root
        server          = /usr/bin/nc
        server_args     = 127.0.0.1 2003
}
                                                                                
